# etc/systemd/systemd-fsck@.service.d/override.conf
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/systemd-fsck %f
StandardOutput=null
StandardError=journal+console
TimeoutSec=0